# config/systemd/minecraft-backup.service
[Unit]
Description=Back up Minecraft Paper server data
Documentation=file:/opt/minecraft-cloud/docs/backup-and-restore.md
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=root
Group=root
Environment=RETENTION_DAYS=7
ExecStart=/opt/minecraft-cloud/scripts/backup.sh
TimeoutStartSec=30min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7